# olegdb/runit/olegdb/run
#!/bin/sh

cd /oleg
/usr/local/bin/olegdb -config /oleg/db.conf
