# net/quassel/runit/quassel/run
#!/bin/sh

quasselcore --configdir=/var/lib/quassel
