# net/gogs/runit/gogs/run
#!/bin/sh

cd /gogs/install
setuser git ./gogs web
