# base/old-alpine/runit/backplane/run
#!/bin/sh

if [ -n "$BACKPLANE_TOKEN" ]
then
	/usr/bin/backplane connect
fi
