# dev/guix/runit/guix-daemon/run
#!/bin/sh

~root/.guix-profile/bin/guix-daemon --build-users-group=guixbuild
